# mitro-core/production/ansible/changepwagent/init/changepwagent.conf
start on runlevel [2345]
stop on runlevel [!2345]
respawn
console log

setuid changepwagent
setgid changepwagent

chdir /home/ubuntu/changepwagent

exec ./changepwagent
