# mitro-core/production/ansible/primary/init/mitrosync.conf
start on runlevel [2345]
stop on runlevel [!2345]
respawn
console log

setuid mitrosync
setgid nogroup

chdir /home/ubuntu/mitrosync/temp
exec java -ea -jar ../mitrosync.jar
